# systemd/hellotalk-process-audio.timer
[Unit]
Description=Daily processing of HelloTalk audio at 12:05

[Timer]
OnCalendar=*-*-* 12:05:00
Persistent=true

[Install]
WantedBy=timers.target
